# vllm/systemd/user/vllm-model@.service
[Unit]
Description=vLLM model service (%i)
After=network-online.target
Wants=network-online.target
PartOf=vllm-stack.target

[Service]
Type=simple
EnvironmentFile=%h/.config/vllm-models/project.env
ExecStart=/usr/bin/env bash -lc '${PROJECT_ROOT}/scripts/run-vllm-server.sh "%i"'
WorkingDirectory=%h
Restart=on-failure
RestartSec=5
TimeoutStartSec=1800
TimeoutStopSec=60
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
