# vllm/systemd/user/vllm-stack.target
[Unit]
Description=vLLM model stack
After=network-online.target
Wants=vllm-model@embeddings.service
Wants=vllm-model@dense.service
Wants=vllm-model@sparse.service
Wants=vllm-model@reranker.service
Wants=vllm-model@coding.service

[Install]
WantedBy=default.target
